# dispatcharr/rootfs/etc/services.d/dispatcharr/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# Ensure DISPATCHARR_ENV is set for aio mode
export DISPATCHARR_ENV="${DISPATCHARR_ENV:-aio}"

# Ingress configuration - service must bind to 0.0.0.0:9191 for Home Assistant Ingress
export INGRESS_PORT=9191
export GUNICORN_PORT=9191
export GUNICORN_BIND="0.0.0.0:9191"

# The actual entrypoint is /app/docker/entrypoint.sh (uses uWSGI + nginx)
# Set NGINX_PORT to match ingress_port in config.yaml (9191)
export NGINX_PORT=9191
export NGINX_LISTEN="0.0.0.0:9191"

# Additional environment variables that might be needed
export PORT=9191
export HOST=0.0.0.0

# Function to patch nginx config to listen on correct port for Ingress
patch_nginx_config() {
    local nginx_conf="${1:-}"
    local patched=0
    
    if [ -z "$nginx_conf" ] || [ ! -f "$nginx_conf" ]; then
        # Try common nginx config locations
        for conf in /etc/nginx/conf.d/default.conf /etc/nginx/sites-enabled/default /app/docker/nginx.conf /etc/nginx/nginx.conf /etc/nginx/sites-available/default; do
            if [ -f "$conf" ]; then
                nginx_conf="$conf"
                break
            fi
        done
    fi
    
    if [ -n "$nginx_conf" ] && [ -f "$nginx_conf" ]; then
        # Check if config already has the correct port
        if grep -q "listen[[:space:]]*0\.0\.0\.0:$NGINX_PORT" "$nginx_conf" 2>/dev/null; then
            echo "Nginx config already set to listen on 0.0.0.0:$NGINX_PORT" >&2
            return 0
        fi
        
        echo "Patching nginx config: $nginx_conf to listen on 0.0.0.0:$NGINX_PORT" >&2
        # Replace any listen directives with our ingress port
        sed -i.bak "s/listen[[:space:]]*[0-9]*;/listen 0.0.0.0:$NGINX_PORT;/g" "$nginx_conf" 2>/dev/null && patched=1
        sed -i.bak "s/listen[[:space:]]*[0-9]*[[:space:]]*default_server;/listen 0.0.0.0:$NGINX_PORT default_server;/g" "$nginx_conf" 2>/dev/null && patched=1
        sed -i.bak "s/listen[[:space:]]*0\.0\.0\.0:[0-9]*/listen 0.0.0.0:$NGINX_PORT/g" "$nginx_conf" 2>/dev/null && patched=1
        sed -i.bak "s/listen[[:space:]]*[0-9]*[[:space:]]*default_server;/listen 0.0.0.0:$NGINX_PORT default_server;/g" "$nginx_conf" 2>/dev/null && patched=1
        rm -f "${nginx_conf}.bak" 2>/dev/null || true
        
        if [ $patched -eq 1 ]; then
            echo "Nginx config patched successfully" >&2
            return 0
        fi
    fi
    
    echo "Warning: Could not find or patch nginx config file" >&2
    return 1
}

# Debug: List all possible entrypoint locations
echo "=== Searching for entrypoint scripts ===" >&2
for ep in /entrypoint.aio.sh /entrypoint.sh /app/docker/entrypoint.sh /docker-entrypoint.sh; do
    if [ -f "$ep" ]; then
        echo "Found: $ep" >&2
        ls -la "$ep" >&2
    fi
done

echo "=== Ingress Configuration ===" >&2
echo "INGRESS_PORT: $INGRESS_PORT" >&2
echo "NGINX_PORT: $NGINX_PORT" >&2
echo "GUNICORN_PORT: $GUNICORN_PORT" >&2
echo "DISPATCHARR_ENV: $DISPATCHARR_ENV" >&2

# Priority order: /app/docker/entrypoint.sh (actual) > entrypoint.aio.sh > entrypoint.sh
if [ -f /app/docker/entrypoint.sh ]; then
    echo "Starting Dispatcharr with /app/docker/entrypoint.sh" >&2
    echo "  Environment: $DISPATCHARR_ENV" >&2
    echo "  NGINX_PORT: $NGINX_PORT (should bind to 0.0.0.0:$NGINX_PORT)" >&2
    
    # Try to patch nginx configs before entrypoint starts (in case they exist)
    patch_nginx_config
    
    # Start entrypoint - it will manage nginx, but we'll need to patch config after nginx starts
    # Use a background process to continuously check and patch nginx config
    (
        while true; do
            sleep 3
            # Check if nginx is running but on wrong port, then patch config and reload
            if command -v nginx >/dev/null 2>&1 && pgrep nginx >/dev/null 2>&1; then
                patch_nginx_config
                # Try to reload nginx if config was patched
                nginx -s reload 2>/dev/null || true
            fi
        done
    ) &
    PATCH_PID=$!
    
    # Start the entrypoint (this will block)
    exec /app/docker/entrypoint.sh
elif [ -f /entrypoint.aio.sh ]; then
    echo "Starting Dispatcharr with /entrypoint.aio.sh" >&2
    echo "  Environment: $DISPATCHARR_ENV" >&2
    echo "  GUNICORN_PORT: $GUNICORN_PORT (should bind to 0.0.0.0:$GUNICORN_PORT)" >&2
    exec /entrypoint.aio.sh
elif [ -f /entrypoint.sh ]; then
    echo "Starting Dispatcharr with /entrypoint.sh" >&2
    echo "  Environment: $DISPATCHARR_ENV" >&2
    echo "  GUNICORN_PORT: $GUNICORN_PORT (should bind to 0.0.0.0:$GUNICORN_PORT)" >&2
    exec /entrypoint.sh
else
    echo "ERROR: No entrypoint script found. Checked:" >&2
    echo "  - /app/docker/entrypoint.sh" >&2
    echo "  - /entrypoint.aio.sh" >&2
    echo "  - /entrypoint.sh" >&2
    echo "Available files in /app/docker/:" >&2
    ls -la /app/docker/ 2>&1 || true
    exit 1
fi
