# dispatcharr/rootfs/etc/services.d/dispatcharr/run
#!/usr/bin/with-contenv bash
set -euo pipefail
exec /entrypoint.sh
