# dispatcharr/rootfs/etc/services.d/dispatcharr/run
#!/usr/bin/with-contenv bash
set -euo pipefail
# Use the upstream entrypoint - this works locally on dev image
exec /entrypoint.sh
